# bandwidth_occupier.timer
[Unit]
Description=Run the Bandwidth Occupier every 60 minutes

[Timer]
OnBootSec=15min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target
